# deploy/nhc.service
[Unit]
Description=NHC — Nethack-like Crawler web server
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=10

# Secrets and tunables live in the override file:
#   /etc/systemd/system/nhc.service.d/override.conf
# Required:  NHC_AUTH_TOKEN
# Optional:  NHC_MAX_SESSIONS (default 8)
#            NHC_GEN_WORKERS  (default 4 — dungeon generation pool)
Environment=NHC_DATA_DIR=/var/nhc
Environment=NHC_MAX_SESSIONS=8
Environment=NHC_GEN_WORKERS=4
Environment=NHC_BIND=0.0.0.0

ExecStartPre=-/usr/bin/docker rm -f nhc
ExecStart=/usr/bin/docker run --name nhc \
    --publish ${NHC_BIND}:8080:8080 \
    --volume /var/nhc:/var/nhc \
    --env NHC_DATA_DIR=${NHC_DATA_DIR} \
    --env NHC_AUTH_TOKEN=${NHC_AUTH_TOKEN} \
    --env NHC_MAX_SESSIONS=${NHC_MAX_SESSIONS} \
    --env NHC_GEN_WORKERS=${NHC_GEN_WORKERS} \
    --env NHC_EXTERNAL_URL=${NHC_EXTERNAL_URL} \
    nhc-web
ExecStop=/usr/bin/docker stop -t 15 nhc
ExecStopPost=-/usr/bin/docker rm -f nhc

[Install]
WantedBy=multi-user.target
